# systemd/kern.service
[Unit]
Description=Kern Smart Resource Manager for Developers
Documentation=https://github.com/yashurade27/kern
After=graphical-session-pre.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/kern enforce
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (optional)
CPUAccounting=true
MemoryAccounting=true

[Install]
WantedBy=graphical-session.target
